# roles/custom_caching_server/templates/cache-server.service.j2
[Unit]
Description= Build and Run Rust Cache server

[Service]
ExecStart=/bin/bash /root/run-cache-server.sh

[Install]
WantedBy=default.target
